# conky/weather/conky.weather.cron
*/20 9-18 * * 1-5 $HOME/.conky/weather/getWeatherData.sh > /dev/null 2>&1
